# packages/python-services/resources/systemd/geth-ropsten.service
[Unit]
Description=Ethereum Go client
OnFailure=unit-status-mail@%n.service

[Service]
ExecStart=/usr/bin/geth --testnet --fast --rpc --rpcaddr localhost --rpcport 8546 --rpcapi "eth,net,web3,admin,personal,txpool" --port 30304
Restart=always
RestartSec=30
Type=simple
User=ed

[Install]
WantedBy=default.target
