# packages/python-services/resources/systemd/geth.service
[Unit]
Description=Ethereum Go client
OnFailure=unit-status-mail@%n.service

[Service]
ExecStart=/usr/bin/geth --fast --rpc --rpcaddr 46.101.132.175 --rpcapi "eth,net,web3,admin,personal,txpool"
Restart=always
RestartSec=30
Type=simple
User=ed

[Install]
WantedBy=default.target
