# packages/python-services/resources/systemd/geth.service
[Unit]
Description=Ethereum Go client
OnFailure=unit-status-mail@%n.service

[Service]
ExecStart=/usr/bin/geth --syncmode "light" --rpc
Restart=always
RestartSec=30
Type=simple
User=ed

[Install]
WantedBy=default.target
