# packages/python-services/resources/systemd/rcpyd.service
[Unit]
Description=Realitio Service Daemon
OnFailure=unit-status-mail@%n.service

[Service]
WorkingDirectory=/home/ed/working/realitycheck-python-services/src/rcpy
ExecStart=/home/ed/venvs/rcpy/bin/python3 manage.py daemon
Restart=always
RestartSec=60
Type=simple
User=ed

[Install]
WantedBy=default.target
